# rootfs/etc/services.d/jenkins/run
#!/usr/bin/with-contenv sh
                                                                                                                 
exec s6-setuidgid jenkins /usr/local/bin/jenkins.sh